# bin/chatserver.service
[Unit]

Description=Django web socket chat demo.

After=network.target

 

[Service]

User=user

Environment="DJANGO_SECRET_KEY=sssssomekey"

Environment="DJANGO_SETTINGS_MODULE=WebSocketServer.settings"

Group=www-data

WorkingDirectory=/var/www/websocket_server/

ExecStart=/home/ubuntu/.virtualenvs/websocketserver/bin /var/www/websocket_server/manage.py run_chat_server

[Install]

WantedBy=multi-user.target
